# proxy/wsproxy.conf
description "node.js websocket proxy for Bedlam"
author      "gom"

start on startup
stop on shutdown
respawn
respawn limit 10 5

script
    export HOME="/root"
    exec /usr/bin/node /home/plamen/BEDLAM_CODER/src/wsproxy.js 2>&1 >> /var/log/wsproxy.log
end script
